# root/etc/services.d/gitea/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc \
    gitea -c "$CONF_FILE" web
